# ansible/playbooks/old/files/puma.service
[Unit]
Description=Puma HTTP Server
After=network.target

[Service]
Type=simple
EnvironmentFile=/home/alexgres48/db_config
User=alexgres48
WorkingDirectory=/home/alexgres48/reddit
ExecStart=/bin/bash -lc 'puma'
Restart=always

[Install]
WantedBy=multi-user.target